# dev-daemon/dev-daemon.service.example
[Unit]
Description=Dev stack (backend + frontend)
After=network.target

[Service]
Type=simple
# Set to your dev repo root (directory containing dev-server, dev-frontend, dev-daemon)
WorkingDirectory=/path/to/dev
ExecStart=/path/to/dev/dev-daemon/start.sh
Restart=on-failure
RestartSec=5

# So the service can run without a full login session
Environment=HOME=%h

[Install]
WantedBy=default.target
